# archlinux/etc/sv/agetty-tty2/run
#!/bin/sh

exec /sbin/agetty -8 38400 tty2 linux
