# archlinux/etc/sv/agetty-tty1/run
#!/bin/sh

exec /sbin/agetty -8 38400 tty1 linux
